# examples/service/dev.gosaid.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>dev.gosaid</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/gosaid</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/tmp/gosaid.out</string>
    <key>StandardErrorPath</key>
    <string>/tmp/gosaid.err</string>
</dict>
</plist>
